# systemd-templates/backup-dir.timer
[Unit]
Description=Backup DIR daily

[Timer]
OnCalendar=daily
Persistent=true
 
[Install]
WantedBy=timers.target
