# systemd-templates/backup-dir.service
[Unit]
Description=Backup DIR

[Service]
Type=oneshot
ExecStart=/home/USER/bin/backup-dir SRCDIR DSTDIR
