# systemd-templates/restic-backup.timer
[Unit]
Description=Regular backup using restic

[Timer]
OnCalendar=daily
Persistent=true
 
[Install]
WantedBy=timers.target
